# scripts/com.lgvs.geo.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd 설정 파일
  
  설치:
    1. 이 파일을 ~/Library/LaunchAgents/ 에 복사
    2. 경로를 실제 프로젝트 경로로 수정
    3. launchctl load ~/Library/LaunchAgents/com.lgvs.geo.daily.plist
  
  제거:
    launchctl unload ~/Library/LaunchAgents/com.lgvs.geo.daily.plist
  
  상태 확인:
    launchctl list | grep lgvs
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lgvs.geo.daily</string>
    
    <key>ProgramArguments</key>
    <array>
        <!-- 경로를 실제 프로젝트 경로로 수정하세요 -->
        <string>/Users/YOUR_USERNAME/Documents/LG_VS/scripts/daily-collect.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <!-- 경로를 실제 프로젝트 경로로 수정하세요 -->
    <string>/Users/YOUR_USERNAME/Documents/LG_VS</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <!-- 경로를 실제 프로젝트 경로로 수정하세요 -->
    <string>/Users/YOUR_USERNAME/Documents/LG_VS/logs/launchd.log</string>
    
    <key>StandardErrorPath</key>
    <!-- 경로를 실제 프로젝트 경로로 수정하세요 -->
    <string>/Users/YOUR_USERNAME/Documents/LG_VS/logs/launchd-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
